# image/live/config/includes.chroot/etc/systemd/system/getty@tty1.service.d/zz-autologin.conf
# live-config's getty generator writes /run/.../live-config_autologin.conf with
# "agetty --autologin user". This image has no "user" account, so that getty exits
# 1, hits the start limit, and leaves the VGA console (tty1) blank. systemd applies
# drop-ins in lexical filename order, so this file must sort after
# live-config_autologin.conf to win - hence the zz- prefix. Do not rename it earlier.
[Service]
ExecStart=
ExecStart=-/sbin/agetty --autologin root --noclear %I $TERM
